# debian/rum/etc/cron.d/rum
#
# Regular cron jobs for the rum package
#
* * * * * /root/bin/export_mariadb55_cdb.pl --quiet
